# django/debian/skol-django.service
[Unit]
Description=SKOL Django Web Application
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=skol
Group=skol
WorkingDirectory=/opt/skol/django
ExecStart=/opt/skol/django-venv/bin/with_skol_django skol-django runserver 127.0.0.1:8000
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=append:/var/log/skol/skol-django.log
StandardError=append:/var/log/skol/skol-django.log

# Environment - load from file for customization
# Config file is in shared location to persist across version upgrades
# The - prefix makes the file optional (won't fail if missing)
EnvironmentFile=-/opt/skol/config/skol-django.env
EnvironmentFile=-/home/skol/.skol_env
Environment=DJANGO_SETTINGS_MODULE=skolweb.settings
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
